# DEBIAN/prerm
#!/bin/bash
if [ -d /run/systemd/system ]; then
    systemctl stop  chonos-log.service
fi
systemctl disable chonos-log.service
systemctl daemon-reload
