# DEBIAN/prerm
#!/bin/bash
systemctl disable chonos-log.service
systemctl stop  chonos-log.service
systemctl daemon-reload
